# Layer_4/com.vantage.gitsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vantage.gitsync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mauriciomeyran/Documents/04-VANTAGE_CV/Layer_4/wrappers/git_sync_wrapper.sh</string>
    </array>

    <!-- Corre cada 6 horas: 09:00 · 15:00 · 21:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/vantage_l4_gitsync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vantage_l4_gitsync_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
